# chilli/scripts/etc/hotplug.d/iface/30-chilli
#!/bin/sh

[ "$ACTION" == "ifup" ] || exit 0

/etc/init.d/chilli restart $INTERFACE

