# systemd/spot.service
[Unit]
Description=Hotspot Server
After=mongod.service
Requires=mongod.service
BindsTo=mongod.service

[Service]
Type=simple
ExecStartPre=/bin/sleep 1
ExecStart=/opt/spot4/exe.linux-x86_64-3.6/spot4.exe
WorkingDirectory=/opt/spot4/exe.linux-x86_64-3.6/
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target

